# supervisor-segnog.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:dragonfly]
command=dragonfly --port 6381 --dbfilename dump --nodf_snapshot_format --dir /data/dragonfly
autostart=true
autorestart=true
startsecs=3
startretries=999
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:nats]
command=nats-server --js --store_dir /data/nats --port 4222
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/supervisor/nats.log
stderr_logfile=/var/log/supervisor/nats.log

[program:memory-service]
# FalkorDB now runs in a separate container (see Dockerfile/docker-compose). Wait
# for it at MEMORY_SERVICE_FALKORDB__URL (set in the container env) before starting.
command=bash -c 'for i in $(seq 1 60); do python -c "import os,socket; from urllib.parse import urlparse; u=urlparse(os.environ.get(\"MEMORY_SERVICE_FALKORDB__URL\",\"redis://falkordb:6379\")); s=socket.create_connection((u.hostname,u.port or 6379),2); s.close()" 2>/dev/null && break; echo "[memory-service] Waiting for external FalkorDB..."; sleep 2; done; echo "[memory-service] FalkorDB ready, starting..."; exec memory-service'
directory=/app
autostart=true
autorestart=true
startsecs=10
startretries=999
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=MEMORY_SERVICE_DRAGONFLY__URL="redis://localhost:6381",MEMORY_SERVICE_NATS__URL="nats://localhost:4222",MEMORY_SERVICE_NATS__ENABLED="true",PYTHONPATH="/pip-install"

[eventlistener:sync-on-exit]
command=/app/sync-data.sh
events=PROCESS_STATE_EXITED
autostart=true
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
